# config/systemd/user/watchdog.timer
[Unit]
Description=Guguto's player (RFID + Spotify) timer

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s